# resources/programfiles/open-onedrive-daemon.service
[Unit]
Description=OpenOneDrive Daemon

[Service]
ExecStart=%h/.local/bin/open-onedrive
Restart=on-failure

[Install]
WantedBy=default.target 